# foxnet-updater-fixup.service
[Unit]
Description=Foxnet update configuration update

[Service]
Type=oneshot
ExecStart=/usr/bin/foxnet-updater-fixup

[Install]
WantedBy=multi-user.target
