# foxnet.service
[Unit]
Description=Foxnet service
After=network-online.target

[Service]
ExecStart=/usr/bin/foxnet-init

[Install]
WantedBy=multi-user.target
