# steamos-devkit-service.service
[Unit]
Description="SteamOS Devkit Service"

[Service]
ExecStart=/usr/bin/python /usr/share/steamos-devkit/steamos-devkit-service.py
User=deck

[Install]
WantedBy=multi-user.target
