# steam-web-debug-portforward.socket
[Socket]
ListenStream=8081
ReusePort=yes

[Install]
WantedBy=sockets.target
