# steamos-devkit-service.service
[Unit]
Description="SteamOS Devkit Service"

[Service]
ExecStart=/usr/bin/steamos-devkit-service
User=deck

[Install]
WantedBy=multi-user.target
